# run/dev/php/etc/supervisor.d/rabbitmq.ini
[program:mail-symfony]
command=/srv/skema/bin/console messenger:consume async_mailer --memory-limit=128M --time-limit=3600 --limit=100
autostart=true
autorestart=true

[program:notification-symfony]
command=/srv/skema/bin/console messenger:consume async_notifier --memory-limit=128M --time-limit=3600 --limit=100
autostart=true
autorestart=true

[program:cvcalculation-symfony]
command=/srv/skema/bin/console messenger:consume async_cvcalculator --memory-limit=128M --time-limit=3600 --limit=100
autostart=true
autorestart=true

[program:multicaster-symfony]
command=/srv/skema/bin/console messenger:consume async_multicaster --memory-limit=128M --time-limit=3600 --limit=100
autostart=true
autorestart=true

[program:admissibility_calculator-symfony]
command=/srv/skema/bin/console messenger:consume async_admissibility_calculator --memory-limit=2048M --time-limit=3600 --limit=100
autostart=true
autorestart=true

[program:eligible_students-symfony]
command=/srv/skema/bin/console messenger:consume async_eligible_students --memory-limit=128M --time-limit=3600 --limit=100
autostart=true
autorestart=true

[program:student_export_list-symfony]
command=/srv/skema/bin/console messenger:consume async_student_export_list --memory-limit=2048M --time-limit=3600 --limit=100
autostart=true
autorestart=true

[program:payment_status-symfony]
command=/srv/skema/bin/console messenger:consume async_payment_status --memory-limit=2048M --time-limit=3600 --limit=100
autostart=true
autorestart=true

[program:oral_test_student_reservation]
command=/srv/skema/bin/console messenger:consume async_oral_test_student_reservation --memory-limit=1024M --time-limit=3600 --limit=50
autostart=true
autorestart=true

[program:written_test_summons_message-symfony]
command=/srv/skema/bin/console messenger:consume async_written_test_summons_message --memory-limit=2048M --time-limit=3600 --limit=100
autostart=true
autorestart=true

[program:written_test_summon_message-symfony]
command=/srv/skema/bin/console messenger:consume async_written_test_summon_message --memory-limit=2048M --time-limit=3600 --limit=100
autostart=true
autorestart=true

[program:approve-all-candidacy-symfony]
command=/srv/skema/bin/console messenger:consume async_approve_all_candidacy --memory-limit=2048M --time-limit=3600 --limit=100
autostart=true
autorestart=true

[program:notifications_admissibility_result-symfony]
command=/srv/skema/bin/console messenger:consume async_notifications_admissibility_result --memory-limit=2048M --time-limit=3600 --limit=100
autostart=true
autorestart=true
